# luci-app-mt5700m/root/etc/init.d/mt5700m-manager
#!/bin/sh /etc/rc.common

START=82
STOP=28
USE_PROCD=1

service_triggers() {
	procd_add_reload_trigger mt5700m
}

start_service() {
	/usr/sbin/mt5700m-manager sync >/dev/null 2>&1 || true
	procd_open_instance
	procd_set_param command /usr/sbin/mt5700m-manager monitor
	procd_set_param respawn 3600 5 5
	procd_set_param stdout 0
	procd_set_param stderr 0
	procd_close_instance
}

stop_service() {
	/usr/sbin/mt5700m-manager down >/dev/null 2>&1 || true
}

reload_service() {
	/usr/sbin/mt5700m-manager sync >/dev/null 2>&1 || true
}

extra_command 'connect' 'Enable and start the MT5700M data connection'
extra_command 'disconnect' 'Disable the MT5700M data connection'
extra_command 'redial' 'Recycle the MT5700M network connection'
extra_command 'status_json' 'Show MT5700M manager status as JSON'

connect() { /usr/sbin/mt5700m-manager connect; }
disconnect() { /usr/sbin/mt5700m-manager disconnect; }
redial() { /usr/sbin/mt5700m-manager redial; }
status_json() { /usr/sbin/mt5700m-manager status-json; }
